# systemd/gcs-video-rtsp.service
[Unit]
Description=GCS RTSP video relay for Wi-Fi AP clients
After=network-online.target drone-hotspot.service gcs-video-udp-relay.service
Wants=network-online.target gcs-video-udp-relay.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/gcs-ap-streaming
ExecStart=/usr/local/bin/start-video-rtsp.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
KillMode=mixed
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
