# MAV_widget/systemd/drone-hotspot.service
[Unit]
Description=Drone WiFi AP (uap0 + hostapd + dnsmasq)
After=NetworkManager.service network-online.target
Wants=network-online.target
Conflicts=shutdown.target
ConditionPathExists=!/var/lib/gcs-ap/manual-off

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=180
TimeoutStopSec=30
ExecStart=/usr/local/bin/start-drone-hotspot.sh
ExecStop=/usr/local/bin/stop-drone-hotspot.sh --no-restore

[Install]
WantedBy=multi-user.target
