# MAV_widget/systemd/gcs-wlan-keepalive.timer
[Unit]
Description=Periodic wlan0 reconnect while AP active

[Timer]
OnBootSec=60
OnUnitActiveSec=30
AccuracySec=10
Unit=gcs-wlan-keepalive.service

[Install]
WantedBy=timers.target
