# MAV_widget/systemd/gcs-ap-default-off.service
[Unit]
Description=GCS AP default off at boot (enable only via tray/toggle-ap)
After=local-fs.target
Before=drone-hotspot.service

[Service]
Type=oneshot
ExecStart=/usr/bin/mkdir -p /var/lib/gcs-ap
ExecStart=/usr/bin/touch /var/lib/gcs-ap/manual-off
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
